# php8/base/os15-php8/etc/supervisor.conf
[inet_http_server]
port=*:9010
username=supervisor_admin
password=

[supervisord]
logfile=/var/log/supervisord/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200
nocleanup=true
user=root

[program:crond]
command=cron -n
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=false

[program:apache]
command=/usr/sbin/httpd -DFOREGROUND
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=false
startsecs=5

[program:php-fpm]
command=/usr/sbin/php-fpm -F -R
process_name=%(program_name)s_%(process_num)02d
numprocs=1
autostart=true
autorestart=false
startsecs=5